# rpg.ini
[program:rpg]
command = /Users/decitre/workspace/ng/django-rpg/run.sh -b localhost:8088
autorestart=true
stdout_logfile = /tmp/gunicorn_supervisor.log 
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8 
